# home/private_dot_config/systemd/user/vdirsyncer.service.d/torpedo-credstore.conf
# Provisioned by chezmoi (torpedo only). Loads the TPM2+host sealed Nextcloud
# CardDAV app password from the user credstore so vdirsyncer can fetch it at
# run-time from $CREDENTIALS_DIRECTORY/vdirsyncer-nextcloud-password.
#
# Provision the credential with:
#   systemd-creds encrypt --user --with-key=tpm2+host \
#     -n vdirsyncer-nextcloud-password - \
#     %h/.config/credstore.encrypted/vdirsyncer-nextcloud-password.cred
[Service]
LoadCredentialEncrypted=vdirsyncer-nextcloud-password:%h/.config/credstore.encrypted/vdirsyncer-nextcloud-password.cred
