# Library/LaunchAgents/name.bapt.mbsync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- http://www.launchd.info/ -->
<!-- https://developer.apple.com/library/archive/documentation/MacOSX/Conceptual/BPSystemStartup/Chapters/CreatingLaunchdJobs.html -->
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>name.bapt.mbsync</string>
    <key>ProgramArguments</key>
    <array>
      <string>/usr/local/bin/mbsync</string>
      <string>-a</string>
      <string>-d</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/baptiste</string>
    <key>StandardOutPath</key>
    <string>/Users/baptiste/Mail/.mbsync-agent.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/baptiste/Mail/.mbsync-agent.log</string>
    <key>StartInterval</key>
    <integer>60</integer>
  </dict>
</plist>
