# home/private_dot_config/systemd/user/mailsync.timer
[Unit]
Description=Mailbox synchronization timer

[Timer]
OnBootSec=2m
OnUnitActiveSec=3m
Unit=mailsync.service

[Install]
WantedBy=timers.target
