# home/private_dot_config/systemd/user/pai-git-sync.timer
[Unit]
Description=PAI ~/.claude git sync — every 4h, 7 days a week

[Timer]
OnCalendar=*-*-* 0/4:00:00 Europe/Paris
RandomizedDelaySec=5min
Persistent=true

[Install]
WantedBy=timers.target
