# home/private_dot_config/systemd/user/pai-git-sync.service
[Unit]
Description=PAI ~/.claude git sync
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/.claude
ExecStart=/usr/bin/bash -c '\
  cd %h/.claude && \
  git add -A && \
  if ! git diff --cached --quiet; then \
    git commit -m "auto: timed sync $(date +%%Y-%%m-%%d\ %%H:%%M)"; \
    git push; \
  fi'
StandardOutput=journal
StandardError=journal
