# home/private_dot_config/systemd/user/pai-git-sync.service
[Unit]
Description=PAI ~/.claude git sync
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/.claude
ExecStart=%h/.local/bin/pai-git-sync
StandardOutput=journal
StandardError=journal
